# .config/systemd/user/pywal_morning.timer
[Unit]
Description=Timer for pywal_morning.service

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
