# .config/systemd/user/pywal_evening.timer
[Unit]
Description=Timer for pywal_evening.service

[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
